# systemd/globalprotect-fix-resolvconf.path
[Unit]
Description=fix /etc/resolv.conf after globalprotect overwrites it

[Path]
PathChanged=/etc/resolv.conf
Unit=globalprotect-fix-resolvconf.service

[Install]
WantedBy=multi-user.target
